# etc/fpnd.service
[Unit]
Description=fpnd control daemon
Requires=zerotier-one.service
After=zerotier-one.service

[Service]
ExecStart=/usr/lib/fpnd/fpnd.py start
ExecStop=/bin/kill -TERM $MAINPID
KillMode=none
PIDFile=/run/fpnd.pid
Type=forking
TimeoutStartSec=5
TimeoutStopSec=7

[Install]
WantedBy=multi-user.target
