# etc/fpnd.openrc
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

CONFIG="/etc/fpnd/fpnd.ini"

description="fpnd - FreePN daemon"
command="/usr/sbin/fpnd.py"
pidfile="/run/${SVCNAME}.pid"
command_args="start"

depend() {
#        need zerotier
        use logger clock
}

start_pre() {
	if [ ! -f ${CONFIG} ]; then
		eerror "Configuration file ${CONFIG} not found"
		return 1
	fi
}
